# hf_deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
user=hfuser

[program:chromadb]
command=chroma run --host 0.0.0.0 --port 8000 --path /tmp/chroma
autostart=true
autorestart=true
stderr_logfile=/tmp/chroma.err.log
stdout_logfile=/tmp/chroma.out.log

[program:mcp]
command=uvicorn main:app --host 0.0.0.0 --port 8080
directory=/app/mcp_server
autostart=true
autorestart=true
stderr_logfile=/tmp/mcp.err.log
stdout_logfile=/tmp/mcp.out.log

[program:backend]
command=uvicorn api.main:app --host 0.0.0.0 --port 7860
directory=/app/backend
environment=CHROMA_HOST="http://localhost:8000",MCP_SERVER_HOST="http://localhost:8080"
autostart=true
autorestart=true
stderr_logfile=/tmp/backend.err.log
stdout_logfile=/tmp/backend.out.log
startsecs=5
